# secret-agent.service
[Unit]
Description=SecretAgent Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/SecretAgent
ExecStart=/usr/local/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
